# ansible/roles/elk_stack/templates/elasticsearch_exporter.service.j2
# Generated by Ansible - elk_stack role
[Unit]
Description=Prometheus Elasticsearch Exporter
Documentation=https://github.com/prometheus-community/elasticsearch_exporter
After=network.target elasticsearch.service
Wants=elasticsearch.service

[Service]
Type=simple
User=elasticsearch_exporter
Group=elasticsearch_exporter
ExecStart=/usr/local/bin/elasticsearch_exporter \
  --es.uri={{ elasticsearch_exporter_uri }} \
  --es.all={{ elasticsearch_exporter_all | lower }} \
  --es.indices={{ elasticsearch_exporter_indices | lower }} \
  --es.shards={{ elasticsearch_exporter_shards | lower }} \
  --es.cluster_settings={{ elasticsearch_exporter_cluster_settings | lower }} \
  --web.listen-address=:{{ elk_prometheus_exporter_port }}

SyslogIdentifier=elasticsearch_exporter
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
